# all_new_2018/linkable_etc_files/mail/etc/systemd/system/timers.target.wants/fetchmail.timer
[Unit]
Description=Run fetchmail once every minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
