# all_new_2018/linkable_etc_files/web/etc/systemd/system/timers.target.wants/pull_website.timer
[Unit]
Description=Pull website repo every minute

[Timer]
OnCalendar=*-*-* *:*:00,15,30,45

[Install]
WantedBy=timers.target
