# all_new_2018/linkable_etc_files/play/etc/systemd/system/timers.target.wants/encrypt_chatlogs.timer
[Unit]
Description=Attempt encryption of old chatlogs once every minute.

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target