# all_new_2018/linkable_etc_files/mail/etc/systemd/system/pingmail.service
[Unit]
Description=Run pingmail check

[Service]
Type=oneshot
User=plom
ExecStart=/bin/sh -c '~/pingmail/pingmail check'
