# all_new_2018/linkable_etc_files/web/etc/systemd/system/plomlombot.service
[Unit]
Description=plomlombot screen

[Service]
Type=simple
User=plom
ExecStart=/bin/sh -c '~/plomlombot_daemon.sh'
Restart=always

[Install]
WantedBy=multi-user.target
