# all_new_2018/linkable_etc_files/web/etc/systemd/system/plomlombot.service
[Unit]
Description=plomlombot screen

[Service]
Type=forking
User=plom
ExecStart=/bin/sh -c 'screen -d -m -S plomlombot ~/start_plomlombot.sh plomlombob "#plomlompomtest"'

[Install]
WantedBy=multi-user.target
