# all_new_2018/linkable_etc_files/play/etc/systemd/system/encrypt_chatlogs.service
[Unit]
Description=Pull website repo
[Service]
Type=oneshot
User=plom
ExecStart=/bin/sh -c '~/encrypter.sh'
